# Hacker_News/config_files/crontab
0 * * * * /home/olives/Hacker_News/venv/bin/python /home/olives/Hacker_News/data.py